# cmder/config/profile.d/SbbCitrix.ps1
function Stop-SbbCitrix {
    Get-Process | Select-Object Id,Path | Where-Object { "$_.Path".ToLower().Contains("citrix") } | Stop-Process -Force
}